# systemd/env.service
[Unit]
Description=Playground service which outputs the process environment
After=network.target

[Service]
User=1000
Group=1000
ExecStartPre=env
ExecStart=sleep 180
Restart=always
RestartSec=5
Environment="PATH=/home/nicoretti/.local/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target


